# dev_script/crontab
0 0 * * * /app/dev_script/robot.sh >> /var/log/login.log 2>&1